# io.github.curtiside.powermate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.github.curtiside.powermate</string>
	<key>Program</key>
	<string>/usr/local/bin/powermate</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<!-- __HOME__ is replaced with $HOME by `make install` (launchd does not
	     expand ~). Without these, all diagnostics vanish into launchd. -->
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/powermate.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/powermate.log</string>
</dict>
</plist>
